# server-health-monitor.service
[Unit]
Description=Server Health Monitor Service
After=network.target
Wants=network.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/server-alerts
ExecStart=/usr/bin/env bash -lc '/opt/server-alerts/.venv/bin/python /opt/server-alerts/monitor.py --run-once'
StandardOutput=journal
StandardError=journal
SyslogIdentifier=server-health-monitor

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/opt/server-alerts
ProtectHome=yes

# Restart on failure
Restart=no

[Install]
WantedBy=multi-user.target